# cookbooks/kibana/files/kibana.conf
[program:kibana]
command=/opt/inmobi/kibana/bin/kibana
stdout_logfile=/var/log/supervisor/kibana.log
redirect_stderr=true
autostart=true
directory=/opt/inmobi/kibana/
user=gridops
